# supervisord.d/jupyterhub.ini
[program:jupyterhub]
command=/usr/bin/python3 /srv/jupyterhub/start_jupyterhub.py --no-ssl --config /srv/jupyterhub/jupyterhub_config.py

